# terraria/service/terraria.socket
[Unit]
Description=Socket for terraria.service
PartOf=terraria.service

[Socket]
ListenFIFO=/opt/terraria/terraria.stdin
