# scripts/systemd/factorio.socket
[Unit]
PartOf=factorio.service

[Socket]
ListenFIFO=/home/cristian/servers/factorio/factorio.stdin
