# scripts/systemd/projectzomboid.socket
[Unit]
PartOf=projectzomboid.service

[Socket]
ListenFIFO=/home/cristian/servers/projectzomboid/projectzomboid.stdin
