# scripts/systemd/terraria.socket
[Unit]
PartOf=terraria.service

[Socket]
ListenFIFO=/home/cristian/servers/terraria/terraria.stdin
