# minecraft/service/minecraft.socket
[Unit]
PartOf=minecraft.service

[Socket]
ListenFIFO=/opt/minecraft/minecraft.stdin
