# minecraft/service/minecraft.socket
[Unit]
PartOf=minecraft.service

[Socket]
ListenFIFO=/home/cristian/servers/minecraft/minecraft.stdin
