# projectzomboid/service/projectzomboid.socket
[Unit]
Description=Socket for projectzomboid.service
PartOf=projectzomboid.service

[Socket]
ListenFIFO=/opt/projectzomboid/projectzomboid.stdin
