# factorio/service/factorio.socket
[Unit]
Description=Socket for factorio.service
PartOf=factorio.service

[Socket]
ListenFIFO=/home/cristian/servers/factorio/factorio.stdin
